# studio/systemd/studio-orchestrator.service
[Unit]
Description=Studio Orchestrator — Agent orchestration kernel
Documentation=https://github.com/learhy/studio
After=network.target

[Service]
# Security model: orchestrator runs as an unprivileged user with restricted filesystem access.
# Worker isolation is provided by bubblewrap/Firecracker at the worker level, not here.
# See docs/security.md for the full threat model.
Type=simple
User=studio
Group=studio
ExecStart=/usr/bin/studio-orchestrator
Restart=on-failure
RestartSec=5s

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/lib/studio /run/studio
PrivateTmp=yes
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
# AF_INET/AF_INET6 required for Ollama Cloud, GitHub API, and remote worker TCP connections
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictRealtime=yes
# Disabled: incompatible with Python runtime (requires W+X memory for bytecode compilation)
MemoryDenyWriteExecute=no
LockPersonality=yes

# Runtime directories
RuntimeDirectory=studio
StateDirectory=studio
WorkingDirectory=/var/lib/studio

# Environment
Environment=STUDIO_SETTINGS_PATH=/etc/studio/settings.json
Environment=STUDIO_ORCH_DB_PATH=/var/lib/studio/state.db
Environment=STUDIO_ORCH_SOCKET_PATH=/run/studio/orchestrator.sock

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=studio-orchestrator

[Install]
WantedBy=multi-user.target
